# roles/mail-agent/templates/mail-agent.timer.j2
[Unit]
Description=Run mail agent every {{ mail_agent_poll_interval }}

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ mail_agent_poll_interval }}
Persistent=true

[Install]
WantedBy=timers.target
