# roles/manga-scheduler/templates/manga-chapter-checker.timer.j2
[Unit]
Description=Check for new manga chapters daily

[Timer]
OnCalendar={{ manga_chapter_checker_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
