# roles/minecraft/templates/mc-fake-server.service.j2
[Unit]
Description=Minecraft fake server (shows start message, triggers real server on connect)
After=network.target docker.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 {{ minecraft_data_dir }}/mc-fake-server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mc-fake-server
User={{ minecraft_user }}
Group=docker
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ minecraft_data_dir }}

[Install]
WantedBy=multi-user.target
