# roles/minecraft/templates/mc-idle-watcher.timer.j2
[Unit]
Description=Run Minecraft idle watcher every 5 minutes

[Timer]
OnCalendar={{ minecraft_idle_check_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
