# roles/minecraft/templates/mc-idle-watcher.service.j2
[Unit]
Description=Minecraft idle watcher (shuts down MC after idle timeout)
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ minecraft_data_dir }}/mc-idle-watcher.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mc-idle-watcher
User={{ minecraft_user }}
Group=docker
Environment=DOCKER_CONFIG={{ minecraft_data_dir }}/.docker
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ minecraft_data_dir }}
