# roles/minecraft/templates/mc-proxy.service.j2
[Unit]
Description=Minecraft TCP proxy (always on :{{ minecraft_port }}, routes to MC or fake server)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 {{ minecraft_data_dir }}/mc-proxy.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mc-proxy
TasksMax=64
MemoryMax=256M
LimitNOFILE=256
User={{ minecraft_user }}
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ minecraft_data_dir }}

[Install]
WantedBy=multi-user.target
